# bootstrap/config/systemd/amazon-cloudwatch-agent.service.d/iam-gate.conf
# Gate the vendor-managed amazon-cloudwatch-agent.service on an attached IAM
# role. ExecCondition non-zero => unit is skipped cleanly (condition-failed,
# not failed). Drop-in only; the vendor unit file is never replaced.
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
ExecCondition=/opt/arm64-bootstrap/bin/cloudwatch-iam-gate
