# bootstrap/config/systemd/docker-cleanup.timer
[Unit]
Description=Daily Docker cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
