# bootstrap/config/systemd/docker.service.d/override.conf
[Unit]
RequiresMountsFor=/data
After=network-online.target
After=data.mount
ConditionPathIsMountPoint=/data

[Service]
ExecStartPre=/usr/bin/test -f /data/.mounted-data-volume
Environment="DOCKER_BUILDKIT=1"
Environment="COMPOSE_DOCKER_CLI_BUILD=1"

[Install]
WantedBy=multi-user.target
